# deploy/systemd/dusheng-dpi.service
[Unit]
Description=DuSheng Panel DPI sidecar
Documentation=https://github.com/SatanDS/DuSheng-Panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dusheng-agent
Group=dusheng-agent
WorkingDirectory=/opt/dusheng-agent
RuntimeDirectory=dusheng-dpi
RuntimeDirectoryMode=0750
Environment=DUSHENG_DPI_ENGINE=auto
Environment=DUSHENG_DPI_MAX_FLOWS=8192
Environment=DUSHENG_DPI_FLOW_TTL=2m
Environment=DUSHENG_DPI_MAX_PACKETS=12
ExecStart=/opt/dusheng-agent/dusheng-dpi -listen unix:/run/dusheng-dpi/dpi.sock -engine ${DUSHENG_DPI_ENGINE} -max-flows ${DUSHENG_DPI_MAX_FLOWS} -flow-ttl ${DUSHENG_DPI_FLOW_TTL} -max-packets ${DUSHENG_DPI_MAX_PACKETS}
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/run/dusheng-dpi

[Install]
WantedBy=multi-user.target
